# systemd/spec-kit-web.service
[Unit]
Description=Spec-Kit Web Interface
After=network.target
StartLimitBurst=3
StartLimitInterval=60

[Service]
Type=simple
User=mrmagic
WorkingDirectory=/home/mrmagic/mcp-server-spec-kit
Environment="PATH=/home/mrmagic/.local/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/bin/python3 /home/mrmagic/mcp-server-spec-kit/app.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
