# systemd/spec-kit-mcp.service
[Unit]
Description=MCP Server for GitHub Spec-Kit
After=network.target
Documentation=https://github.com/github/spec-kit

[Service]
Type=simple
User=mrmagic
WorkingDirectory=/opt/mcp
Environment="PATH=/home/mrmagic/.local/bin:/usr/local/bin:/usr/bin:/bin"
# MCP servers run via stdio transport when called by agents
# This service keeps the server script available but doesn't run it as a daemon
ExecStart=/bin/sleep infinity
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10
StandardInput=null
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

